# .docker/app/supervisor.d/mail-worker.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --queue=mails --sleep=3 --tries=3 --max-time=3600
directory=/var/www/html/my_app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/my_app/storage/logs/worker.log
stopwaitsecs=3600
